# config/systemd-user/dws-task-monitor.service
[Unit]
Description=Dev Workspace Task Monitor
After=network-online.target dws-sessions-init.service
Wants=network-online.target dws-sessions-init.service
ConditionPathExists=%h/bin/task-monitor.sh

[Service]
Type=simple
ExecStart=/usr/bin/bash %h/bin/task-monitor.sh
Restart=on-failure
RestartSec=10
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
WorkingDirectory=%h

[Install]
WantedBy=default.target
