# mac-setup/com.wrkflo.chrome-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: auto-start the Chrome CDP bridge on every login so codex on the
  dev-workspace VM can always reach the Mac's Chrome over Tailscale.

  Install with:
    cp mac-setup/com.wrkflo.chrome-cdp.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.wrkflo.chrome-cdp.plist

  Uninstall with:
    launchctl unload -w ~/Library/LaunchAgents/com.wrkflo.chrome-cdp.plist
    rm ~/Library/LaunchAgents/com.wrkflo.chrome-cdp.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wrkflo.chrome-cdp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "$HOME/dev-workspace/mac-setup/chrome-cdp.sh"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <!-- Only restart if the socat bridge dies. Chrome itself we don't babysit. -->
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/chrome-cdp.launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/chrome-cdp.launchagent.err</string>
</dict>
</plist>
